# services/post-scheduler.service
[Unit]
Description=X Post Scheduler
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/x-comment-agent
EnvironmentFile=/root/x-comment-agent/.env
ExecStart=/usr/bin/node /root/x-comment-agent/post-scheduler.js
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
